# services/adsb-backend.service
[Unit]
Description=ADS-B Dashboard Backend API
After=network.target dump1090-fa.service
Wants=dump1090-fa.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/adsb-dashboard/backend
ExecStart=/usr/bin/python3 /opt/adsb-dashboard/backend/adsb_backend.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adsb-backend
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
