# services/adsb-dashboard.service
[Unit]
Description=ADS-B Kivy Touchscreen Dashboard
After=graphical.target adsb-backend.service
Wants=adsb-backend.service

[Service]
Type=simple
User=pi
Group=pi
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /opt/adsb-dashboard/dashboard/main.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adsb-dashboard

[Install]
WantedBy=graphical.target
